# src/scripts/systemd/axlelore.service
[Unit]
Description=AxleLore Vehicle Assistant
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=axlelore
WorkingDirectory=/opt/axlelore
Environment="PYTHONPATH=/opt/axlelore/src"
Environment="PATH=/opt/axlelore/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/axlelore/.venv/bin/python -m uvicorn backend.main:app --host 0.0.0.0 --port 8000 --log-level warning
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Hardening
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/axlelore/data /var/lib/axlelore
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
